# bin/plists/com.mrflip.noatime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- sudo cp ~/bin/plists/com.mrflip.noatime.plist /Library/LaunchDaemons -->
<!-- sudo cp ~/bin/plists/remount_noatime          /usr/local/bin         -->
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.mrflip.noatime</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/remount_noatime</string>
        </array>
        <key>RunAtLoad</key>
        <true />
    </dict>
</plist>
